# templates/usr/lib/systemd/system/bind-rndc-key.service.j2
[Unit]
Description=Generate rndc key for BIND (DNS)

[Service]
Type=oneshot
ExecStart={{ bind.serverroot }}/libexec/generate-rndc-key

